# recipes-connectivity/dhcp/files/dhcpd-wlan0.service
[Unit]
Description=DHCPv4 Server Daemon
Documentation=man:dhcpd(8) man:dhcpd.conf(5)

[Service]
PIDFile=/run/dhcpd-wlan0.pid
ExecStartPre=/bin/touch /var/lib/dhcp/dhcpd.leases
ExecStart=/usr/sbin/dhcpd -f -cf /etc/dhcpd-wlan0.conf -pf /run/dhcpd-wlan0.pid $DHCPDARGS -q wlan0

[Install]
WantedBy=multi-user.target
